# cron.hourly/aptitude-update
#!/bin/sh
aptitude update > /dev/null 2>&1
exit 0
