# cron.hourly/aptitude-update
#!/bin/bash

aptitude update > /dev/null
