# crontabs/failo
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin

@reboot screen -d -m
