# crontabs/leech
SHELL=/bin/bash
LANGUAGE="en_US:en"
LANG="en_US.utf8"
PATH=/home/v-perlbrew/perl5/perlbrew/bin:/home/v-perlbrew/perl5/perlbrew/perls/current/bin:/usr/local/bin:/usr/bin:/bin

# Backup queued torrents
*/5 * * * * rsync -aq --progress ~/torrent/queue/ ~/torrent/queue-powerseeding/ ~/torrent/queue-archive/

# Start the rtorrent screen on system startup
@reboot screen -U -d -m

# Whine if rtorrent has died
*/10 * * * * (ps aux | grep '^leech.* rtorrent' | grep -qv grep || (echo 'rtorrent not running' && exit 1))
