# crontabs/smolder
SHELL=/bin/bash
LANGUAGE="en_US:en"
LANG="en_US.utf8"
PATH=/home/v-perlbrew/perl5/perlbrew/bin:/home/v-perlbrew/perl5/perlbrew/perls/current/bin:/usr/local/bin:/usr/bin:/bin

HOSTNAME="localhost"
@reboot screen -d -m
