# profile.d/v-90-umask.sh
#!/bin/sh

# This is the Debian default
umask 022
