# profile.d/v-99-failo-wisdom.sh
#!/bin/sh

# spread failo's wisdom to interactive shells
if [[ $- == *i* ]] ; then
    failo-wisdom
fi
