# cron.daily/aptitude-clean
#!/bin/bash

aptitude clean > /dev/null
