# cron.daily/etc-git-permissions
#!/bin/sh

chown -R root:adm /etc/.git
chmod -R g+rw /etc/.git
