# cron.d/aptitude-update
23 08 * * *    root    aptitude update > /dev/null 2>&1 ; true
