# cron.d/aptitude-update
23 08 * * *    root    nice -n 15 aptitude update > /dev/null 2>&1 ; true
