# cron.d/disable-ipv6
@reboot root /etc/network/if-pre-up.d/disable_ipv6.sh
