# cron.d/failo-wisdom
PATH=/usr/bin:/usr/local/bin
4 * * * * failo nice -n 15 ionice -c 2 -n 6 /etc/bin/cron/failo-wisdom-update
