# cron.d/mailman-restart
# Deliever mailman queue periodically
*/15 * * * *     root nice -n 15 ionice -c 2 -n 6 /etc/bin/cron/mailman-restart
