# cron.d/munin-munin_stats_plugins
SHELL=/bin/bash

MUNIN_STATS_PLUGINS_CACHE=/var/lib/munin/plugin-state/munin_stats_plugins-cache

0 6 * * * root cd /etc/munin/plugins && (for i in *; do (time -p sudo nice -n 10 munin-run $i) 2>&1 | grep ^real | sed 's/^real //' | tr '\n' ' ' && echo $i; done) > ${MUNIN_STATS_PLUGINS_CACHE}.tmp 2>&1 && mv ${MUNIN_STATS_PLUGINS_CACHE}.tmp $MUNIN_STATS_PLUGINS_CACHE
