# cron.d/munin-cpanp_o
0 * * * * root [ -O /tmp/munin-cpanp_o.cache ] && cpanp -o | egrep '^    [0-9]+' > /tmp/munin-cpanp_o.cache
