# cron.d/openstreetmap-github-mirror
# OpenStreetMap mirrors

# Git repositories
*/5 * * * *  osm-mirror nice -n 15 ionice -c 2 -n 6 ~/g/openstreetmap-mirror/git-mirror.sh openstreetmap-website
*/5 * * * *  osm-mirror nice -n 15 ionice -c 2 -n 6 ~/g/openstreetmap-mirror/git-mirror.sh merkaartor

# JOSM is a special case
*/15 * * * * osm-mirror nice -n 15 ionice -c 2 -n 6 ~/g/openstreetmap-mirror/josm-mirror.sh

# git-svn mirrors
*/45 * * * * osm-mirror nice -n 15 ionice -c 2 -n 6 ~/g/openstreetmap-mirror/svn-mirror.sh potlatch2
* */2 * * *  osm-mirror nice -n 15 ionice -c 2 -n 6 ~/g/openstreetmap-mirror/svn-mirror.sh tirex
* */4 * * *  osm-mirror nice -n 15 ionice -c 2 -n 6 ~/g/openstreetmap-mirror/svn-mirror.sh osm2pgsql
* */4 * * *  osm-mirror nice -n 15 ionice -c 2 -n 6 ~/g/openstreetmap-mirror/svn-mirror.sh potlatch
* */4 * * *  osm-mirror nice -n 15 ionice -c 2 -n 6 ~/g/openstreetmap-mirror/svn-mirror.sh mapnik
* */4 * * *  osm-mirror nice -n 15 ionice -c 2 -n 6 ~/g/openstreetmap-mirror/svn-mirror.sh shp2osm  
* */4 * * *  osm-mirror nice -n 15 ionice -c 2 -n 6 ~/g/openstreetmap-mirror/svn-mirror.sh mod_tile 
* */4 * * *  osm-mirror nice -n 15 ionice -c 2 -n 6 ~/g/openstreetmap-mirror/svn-mirror.sh cgimap
* */12 * * * osm-mirror nice -n 15 ionice -c 2 -n 6 ~/g/openstreetmap-mirror/svn-mirror.sh gosmore
