# cron.d/reniced
# Source filter reniced so that it lists cmd instead of comm,
# i.e. gets the full argument list too

*/10 * * * *    root perl -0777 -pe 's[sub generate_ps_command.*?\{.*?\n\}][sub generate_ps_command { "ps H -eo lwp,cmd" }]s;' < /usr/bin/reniced | perl
