# cron.d/twitter-twatterhose
SHELL=/bin/bash
PATH=/bin:/usr/bin:/usr/local/bin

*/35 * * * * avar (sleep $(($RANDOM \% 600))) && /etc/bin/cron/twatterhose-twat
@daily       avar rm ~/g/bot-twatterhose/twatterhose.brn
