# cron.d/v-awstats-make-configs
PATH=/home/v-perlbrew/perl5/perlbrew/bin:/home/v-perlbrew/perl5/perlbrew/perls/current/bin:/bin:/usr/bin

@weekly root /etc/bin/cron/awstats-make-configs
