# cron.d/v-etc-git-chown-and-chmod
PATH=/bin
* * * * * root chown root:root /etc/sudoers && chmod 0440 /etc/sudoers
