# cron.d/v-etc-git-chown-and-chmod
@hourly root for repo in $(find /etc -name '.git'); do chown -R root:adm $repo && chmod -R u+rw,g+rw,o-w $repo; done
