# cron.d/v-evalserver
PATH=/home/v-perlbrew/perl5/perlbrew/bin:/home/v-perlbrew/perl5/perlbrew/perls/current/bin:/usr/bin

@reboot root evalserver -d -u evalserver
