# cron.d/v-freenet
@daily root cp ~freenet/freenet/freenet.ini /etc/freenet/freenet.ini
