# cron.d/v-github-backup
PATH=/home/v-perlbrew/perl5/perlbrew/bin:/home/v-perlbrew/perl5/perlbrew/perls/current/bin:/bin:/usr/bin

@daily root /etc/github-backup/github-backup
