# cron.d/v-leech-chown-and-chmod
18 */4 * * * root (chgrp -R leech /home/leech/torrent/complete/* && chmod -R g+rw /home/leech/torrent/complete/*)
