# cron.d/v-mailman-restart
@daily root /etc/init.d/mailman restart > /dev/null
