# cron.d/v-mail-restart
PATH=/home/v-perlbrew/perl5/perlbrew/bin:/home/v-perlbrew/perl5/perlbrew/perls/current/bin:/bin:/usr/bin:/usr/sbin:/bin
@hourly root nice -n 10 cronjob -E  -j v-mail-restart -c '/etc/bin/cron/mail-restart'
