# cron.d/v-munin
*/5  * * * *     munin nice -n 15 ionice -c 2 -n 6 /etc/bin/cron/munin-cron
*/15 * * * *     munin nice -n 15 ionice -c 2 -n 6 /etc/bin/cron/munin-cron-ui

