# cron.d/v-munin-dirs
PATH=/home/v-perlbrew/perl5/perlbrew/bin:/home/v-perlbrew/perl5/perlbrew/perls/current/bin:/bin:/usr/bin:/usr/sbin

10 5   * * * root /etc/bin/cron/munin-dirs dirs_root-cache '/[a-oq-z]*'
41 5   * * * root /etc/bin/cron/munin-dirs homedirs-cache '/home/*'
20 5   * * * root /etc/bin/cron/munin-dirs dirs_var_www-cache '/var/www/*'
10 5   * * * root /etc/bin/cron/munin-dirs dirs_var_backup_github-cache '/var/backup/github/*'
8  */6 * * * root /etc/bin/cron/munin-dirs dirs_var_backup_mysql-cache '/var/backup/mysql/*'
15 5   * * * root /etc/bin/cron/munin-dirs dirs_var_lib_openstreetmap_mirror-cache '' 'du -sk /home/osm-mirror/mirrors/* | sed "s|/home/osm-mirror/mirrors|/var/lib/openstreetmap/mirror|"'
