# cron.d/v-mysql-git-backup
@daily root user=root password=$(cat /etc/mysql/passwd/root) /etc/bin/cron/mysqldump-to-git-all
