# cron.d/v-openvpn
@daily root /etc/bin/cron/openvpn-key-package.sh
