# cron.d/v-openstreetmap-github-mirror-git-repack-and-gc
@daily osm-mirror /etc/bin/cron/git-repack-and-gc-dir /home/osm-mirror/mirrors
