# cron.d/v-osm-nix-is
PATH=/home/v-perlbrew/perl5/perlbrew/bin:/home/v-perlbrew/perl5/perlbrew/perls/current/bin:/bin:/usr/bin:/usr/sbin

##
## avar's daily mkgmap generation
##

## osm.nix.is

# Download the map data from geofabrik.de
0 7      * * *  avar perl /var/www/osm.nix.is/script/download-osm-with-retry --out-dir /var/www/osm.nix.is/dump --fallback http://download.geofabrik.de/osm/europe/iceland.osm.bz2

# Database import. There's a munin script that scrapes this data from
# the database

#5  7     * * *  root service postgresql-8.4 start > /dev/null
10 7     * * *  avar /var/www/osm.nix.is/script/daily-database-import.pl
#50 7     * * *  root service postgresql-8.4 stop > /dev/null

# Generate Garmin map / osmdiff
0  8     * * *  avar /var/www/osm.nix.is/script/generate-garmin-img-from-osm.sh
30 8     * * *  avar LC_ALL=en_US.UTF8 perl /var/www/osm.nix.is/script/generate-daily-osmdiff.pl --date=$(date --iso-8601) 2>&1 | ack -v ^log4j
