# cron.d/v-postfix-reload
PATH=/usr/sbin
@daily root service postfix reload
@daily root postmap /etc/postfix/virtual
@daily root postmap /etc/postfix/transport
