# cron.d/v-postfix-reload
PATH=/usr/sbin:/usr/bin:/bin
@daily root service postfix reload | grep -v '^Reloading' && postmap /etc/postfix/virtual && postmap /etc/postfix/transport
