# cron.d/v-postfix-reload
PATH=/usr/sbin:/usr/bin:/bin
@daily root postmap /etc/postfix/virtual && postmap /etc/postfix/transport
