# cron.d/v-revert-debian-motd-changes
@daily root cd /etc && git checkout /etc/motd
