# cron.d/v-revert-debian-motd-changes
@daily root git --git-dir=/etc/.git checkout /etc/motd
