# cron.d/v-sql-git-backup
# MySQL dumps
0 */6 * * * root user=root password=$(cat /etc/mysql/passwd/root) /etc/bin/cron/mysqldump-to-git-all
5 */6 * * * root /etc/bin/cron/git-repack-and-gc-dir /var/backup/mysql

# SQLite dumps
2 */6 * * * smolder type=sqlite shortname=smolder database=/home/smolder/smolder/smolder.sqlite /etc/bin/cron/sqldump-to-git
7 */6 * * * smolder /etc/bin/cron/git-repack-and-gc-dir /var/backup/sqlite
