# cron.d/v-test-etc
PATH=/home/v-perlbrew/perl5/perlbrew/bin:/home/v-perlbrew/perl5/perlbrew/perls/current/bin:/bin:/usr/bin

# Run tests or reboot
@reboot nobody (sleep 300 && perl -e 'my $out = qx[prove -vr /etc/t 2>&1]; if ($?) { print $out; exit 1 }')

# Regenerate tests
0  */4      * * *  root    find /etc -name 'gen.PL' -exec perl {} \;

# Run tests
2  */4      * * *  nobody  perl -e 'my $out = qx[prove -vr    /etc/bind/t 2>&1]; if ($?) { print $out; exit 1 }'
6  */4      * * *  nobody  perl -e 'my $out = qx[prove -vr /etc/apache2/t 2>&1]; if ($?) { print $out; exit 1 }'
10 */4      * * *  nobody  perl -e 'my $out = qx[prove -vr   /etc/mysql/t 2>&1]; if ($?) { print $out; exit 1 }'
14 */4      * * *  nobody  perl -e 'my $out = qx[prove -vr     /etc/tor/t 2>&1]; if ($?) { print $out; exit 1 }'
16 */4      * * *  nobody  perl -e 'my $out = qx[prove -vr  /etc/t/system 2>&1]; if ($?) { print $out; exit 1 }'

