# cron.d/v-twitter-twatterhose
SHELL=/bin/bash
PATH=/home/v-perlbrew/perl5/perlbrew/bin:/home/v-perlbrew/perl5/perlbrew/perls/current/bin:/bin:/usr/bin

*/35 * * * * avar (sleep $(($RANDOM \% 600))) && /etc/bin/cron/twatterhose-twat
@daily       avar rm ~/g/bot-twatterhose/twatterhose.brn
