# bin/cron/crontabs
#!/bin/sh
CRONTABS="avar hinrik leech failo smolder"

for crontab in $CRONTABS; do
    # Skip cron(1) boilerplate
    perl -ne 'next LINE if $. <= 3; print' /var/spool/cron/crontabs/$crontab > /etc/crontabs/$crontab
done
