# bin/cron/crontabs
#!/bin/sh
CRONTABS="avar hinrik leech failo smolder v-perlbrew"

for crontab in $CRONTABS; do
    if test -f $crontab
    then
        # Skip cron(1) boilerplate
        perl -ne 'next LINE if $. <= 3; print' /var/spool/cron/crontabs/$crontab > /etc/crontabs/$crontab
    fi
done
